# systemd/wirescaled.socket
[Unit]
Description=Wirescale socket

[Socket]
ListenStream=/run/wirescale/wirescaled.sock
FlushPending=yes

[Install]
WantedBy=sockets.target
